# recipes-pibert/pibertservice/pibertservice/pibert.service
[Unit]
Description=Test

[Service]
Type=simple
ExecStart=/usr/bin/echo "hello, world!"

[Install]
WantedBy=multi-user.target
